# conf.d/werkzeug.conf
[program:werkzeug]
command=.../virtualenvs/.../bin/python .../quoteserver/quotes.py
directory=.../quoteserver
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/werkzeug.log
